# setup/etc/systemd/system/kresd@.service.d/override.conf
[Service]
EnvironmentFile=-/root/antizapret/setup
LimitNPROC=infinity
LimitNOFILE=524288